# infra/systemd/ibeam.service
[Unit]
Description=IBeam (IBKR Client Portal Gateway authenticator) for bull-call
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
Restart=always
RestartSec=15
TimeoutStartSec=300
EnvironmentFile=/etc/bull-call/ibeam.env

# Pull the IBeam image and our PYTHONPATH for the SSM provider module.
# IMAGE PINNING: voyz/ibeam:0.5.12 + manifest sha256 digest. The tag drives
# readability; the digest defeats tag-overwrite supply-chain attacks. To
# ratchet, fetch the latest release from
# https://github.com/Voyz/ibeam/releases and pull its manifest digest from
# https://hub.docker.com/v2/repositories/voyz/ibeam/tags/<tag>/.
ExecStartPre=-/usr/bin/docker rm -f ibeam
ExecStartPre=/usr/bin/docker pull voyz/ibeam:0.5.12@sha256:7ca5cf13ec8f0fc16aa8f10fd2613bff9add111a431581eb7818ea0ceeb7cd31

# Mount /opt/bull-call/src so IBeam can import bull_call.ibeam_ssm_provider.
# Pass through SSM_PREFIX + AWS_REGION (NOT credentials).
ExecStart=/usr/bin/docker run --rm --name ibeam \
    --network host \
    -e SSM_PREFIX \
    -e AWS_REGION \
    -e IBEAM_GATEWAY_BASE_URL \
    -e IBEAM_LOG_LEVEL \
    -e IBEAM_SECRETS_SOURCE \
    -e PYTHONPATH=/opt/bull-call/src \
    -v /opt/bull-call/src:/opt/bull-call/src:ro \
    voyz/ibeam:0.5.12@sha256:7ca5cf13ec8f0fc16aa8f10fd2613bff9add111a431581eb7818ea0ceeb7cd31

ExecStop=/usr/bin/docker stop ibeam

# Mirror stdout/stderr to a dedicated file for CloudWatch shipping.
StandardOutput=append:/var/log/bull-call/ibeam.log
StandardError=append:/var/log/bull-call/ibeam.log

[Install]
WantedBy=multi-user.target
